# ansible/templates/supervisor.ini
[program:{{ app_name }}]
command=/usr/local/bin/uwsgi
  --chdir={{ webapps_dir }}/{{ app_name }}/src/{{ app_base }}
  --socket /tmp/{{ app_name }}.sock
  --logto {{ webapps_dir }}/{{ app_name }}/log/application.log
  --home {{ webapps_dir }}/{{ app_name }}/venv 
  --wsgi-file {{ webapps_dir }}/{{ app_name }}/src/{{ app_base }}/{{ wsgi_file }}
  --callable {{ wsgi_callable }}
  --max-requests 1000
  --harakiri 600
  --master
  --processes 5
  --chmod
autostart=true
autorestart=true